# ansible/roles/pmkdata/templates/pmkdata.socket
[Unit]
Description=Gunicorn socket for pmkdata

[Socket]
ListenStream=/run/pmkdata.sock

[Install]
WantedBy=sockets.target
